# project-lifecycle-manager.service
[Unit]
Description=Project Lifecycle Manager - Docusaurus Wiki
After=network.target

[Service]
Type=simple
User=kau005
WorkingDirectory=/home/kau005/project-lifecycle-manager
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
Environment=NODE_ENV=production
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=project-lifecycle-manager

[Install]
WantedBy=multi-user.target